# recipes-core/resize-rootfs/files/resize-rootfs.service
[Unit]
Description=Resize rootfs

[Service]
Type=oneshot
ExecStart=/usr/bin/resize-rootfs start

[Install]
WantedBy=multi-user.target
